# analytics/assembly/docker/dlabs/supervisord-spark-driver.conf
[program:spark-driver]
command=/usr/local/spark/bin/spark-submit --class com.dematic.labs.analytics.ingestion.sparks.drivers.EventCount --master yarn-cluster --driver-memory 1g --executor-memory 1g --executor-cores 1 /usr/local/spark/lib/spark-drivers-1.0-SNAPSHOT.jar https://kinesis.us-west-2.amazonaws.com eventstream