# analytics/assembly/docker/dlabs/supervisord-spark.conf
[program:spark]
command=/etc/bootstrap.sh